# systemd/pineapple-update.service
[Unit]
Description=Update pineapple EA

[Service]
ExecStart=/home/deck/Applications/pineapple-update
