# systemd/pineapple-update.service
[Unit]
Description=Update pineapple EA

[Service]
ExecStart=%h/.local/bin/pineapple-update
